# runpod_worker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:comfyui]
command=python /ComfyUI/main.py --listen 0.0.0.0 --port 8188
directory=/ComfyUI
autostart=true
autorestart=true
stderr_logfile=/var/log/comfyui.err.log
stdout_logfile=/var/log/comfyui.out.log

[program:runpod-worker]
command=python -u /src/handler.py
directory=/
autostart=true
autorestart=true
stderr_logfile=/var/log/worker.err.log
stdout_logfile=/var/log/worker.out.log
